# scripts/hermes-hub-linux-update.service
[Unit]
Description=Update Hermes Hub Linux gateway helper
Documentation=https://github.com/JackoPeru/app-interazione-nemoclaw

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hermes-hub-linux-update --restart
