# scripts/hermes-hub-linux-update.timer
[Unit]
Description=Daily Hermes Hub Linux gateway helper update check

[Timer]
OnBootSec=2min
OnCalendar=*-*-* 04:20:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
